# deploy/systemd/llm-switch.service
[Unit]
Description=AI Proxy Server for Dynamic Model Switching
After=network.target

[Service]
Type=simple

ExecStart={{PROJECT_ROOT}}/.venv/bin/python {{PROJECT_ROOT}}/src/systemd_llm_switch/main.py

WorkingDirectory={{PROJECT_ROOT}}/src/systemd_llm_switch/

Restart=always
RestartSec=5

Environment="PATH={{PROJECT_ROOT}}/.venv/bin"

[Install]
WantedBy=default.target